# etc/s6/services/php-fpm8/run
#!/bin/execlineb -P
s6-envdir /run/env
/usr/sbin/php-fpm8
