# etc/s6/services/php-fpm8/run
#!/usr/bin/execlineb -P
/usr/sbin/php-fpm8
