# etc/s6/services/php-fpm82/run
#!/bin/execlineb -P
/usr/sbin/php-fpm82
