# etc/s6/services/php-fpm84/run
#!/bin/execlineb -P
/usr/sbin/php-fpm84
