# etc/s6/services/nginx/run
#!/usr/bin/execlineb -P
if { test -S /var/run/php-fpm.sock }
/usr/sbin/nginx
