# etc/s6/services/php-fpm83/run
#!/bin/execlineb -P
/usr/sbin/php-fpm83
