# etc/init.d/rc.local
#!/bin/execlineb -P
foreground { cp -r /etc/s6/services /run }
/bin/s6-svscan /run/services