# etc/s6/services/php-fpm85/run
#!/bin/execlineb -P
/usr/sbin/php-fpm85
