# etc/s6/services/php-fpm81/run
#!/bin/execlineb -P
/usr/sbin/php-fpm81
